# files/sticky-world-writeable.cron
5 4 * * * root /usr/share/security_baseline/bin/sticky-world-writable.sh > /dev/null 2>&1
