# files/system-file-permissions.cron
3 5 * * * root /usr/local/sin/system-file-permissions.sh > /dev/null 2>&1
