# files/unowned-files.cron
5 23 * * * root /usr/local/sbin/unowned_files.sh > /dev/null 2>&1
