# files/world-writeable-files.cron
5 1 * * * root /usr/share/security_baseline/bin/world-writeable-files.sh > /dev/null 2>&1
