# files/world-writeable-files.cron
5 1 * * * root /usr/local/sbin/world-writeable-files.sh > /dev/null 2>&1
